# examples/upstart.conf
description "bohr"

exec 2>&1

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

console log

setuid bohr
setgid bohr

chdir /var/run/bohr
script
  /usr/local/bin/bohr --config /etc/bohr/bohr.yml --config /etc/bohr/conf.d --periodic >> /var/log/bohr/bohr.log 2>&1
end script  
